# templates/cuckoo_super.conf.j2
[program:cuckoo]
command = {{ cuckoo_root }}/web/gunicorn_start.sh
user = cuckoo
stdout_logfile = /var/log/supervisor/cuckoo.log
redirect_stderr = true
